# scripts/ssmspl-host-actions.service
[Unit]
Description=SSMSPL SuperAdmin host-action daemon
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ssmspl-host-action-daemon.sh
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
# Hardening — daemon needs root for systemctl/docker, but lock down everything else.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ssmspl-host-actions /var/log /tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
